# profile.d/efl.sh
## Elementary
export ELM_PROFILE="wearable"
export ELM_DISPLAY="wl"

## Ecore
export ECORE_IMF_INPUT_PANEL_ENABLED=1

## Evas
export EVAS_FONT_DPI=72
export EVAS_GL_NO_BLACKLIST=1
export EVAS_GL_NOPRELOAD=1

## Eina
##export EINA_LOG_DLOG_ENABLE=1
export EINA_LOG_SYSLOG_ENABLE=0
export EINA_LOG_ABORT=0
export EINA_LOG_ABORT_LEVEL=0

## Wayland
export WAYLAND_DISPLAY="wayland-0"

## Ector
export ECTOR_BACKEND=default
