# profile.d/efl.sh
## Elementary
export ELM_PROFILE="wearable"
export ELM_DISPLAY="wl"
export ELM_ACCEL="gl"

## Ecore
export ECORE_IMF_INPUT_PANEL_ENABLED=1

## Evas
export EVAS_FONT_DPI=72
export EVAS_GL_NO_BLACKLIST=1
export EVAS_GL_NOPRELOAD=1

## Wayland
export WAYLAND_DISPLAY="wayland-0"

## Ector
export ECTOR_BACKEND=default
