# etc/profile.d/efl.sh
## Elementary
export ELM_MODULES="ctxpopup_copypasteUI>entry/api:datetime_input_ctxpopup>datetime/api"
export ELM_MODULES="$ELM_MODULES:prefs>prefs_iface"

## Ecore
export ECORE_IMF_INPUT_PANEL_ENABLED=1

## Evas
export EVAS_FONT_DPI=72
export EVAS_GL_NO_BLACKLIST=1
