# etc/profile.d/elm.sh
export ELM_MODULES="ctxpopup_copypasteUI>entry/api:datetime_input_ctxpopup>datetime/api"
export ELM_MODULES="$ELM_MODULES:prefs>prefs_iface"
